# .ref/archive-2022/.ref-cdk/clickhouse/ami/roles/clickhouse/files/clickhouse-server.path
[Unit]
Description=Monitor the file for presence. As soon as file appears (volume is mounted) - start clickhouse.

[Path]
PathExists=/var/lib/clickhouse/.ec2-vol
Unit=clickhouse-server.service

[Install]
WantedBy=multi-user.target
