# .ref/archive-2022/.ref-cdk/clickhouse/ami/roles/__aws-cfn-helpers/files/cfn-hup.service
[Unit]
Description=cfn-hup daemon

[Service]
Type=simple
ExecStart=/usr/local/bin/cfn-hup
Restart=always

[Install]
WantedBy=multi-user.target
